# airootfs/etc/xdg/autostart/Calamares.desktop
[Desktop Entry]
Type=Application
Name=Calamares
Exec=sh -c "sudo -E calamares -D8"
