# airootfs/etc/xdg/autostart/Calamares.desktop
[Desktop Entry]
Type=Application
Name=Calamares
Exec=sudo -E calamares
